# Recources/backup/internet-connection.service
[Unit]
Description=Internet Connection Setup

[Service]
ExecStart=/home/nikla/internet-connection.sh

[Install]
WantedBy=multi-user.target
